# server/ble/sentryusb-telemetry.service
[Unit]
# Naming note: this unit is `sentryusb-telemetry.service` but the
# binary it runs is `sentryusb-tesla-telemetry` (the crate name in
# Cargo.toml is `sentryusb-tesla-telemetry` too). The asymmetry is
# historical — the unit was named before the "tesla-" prefix was
# added — and a rename now would break upgrade paths on every
# existing install (systemd would lose track of the old unit). Use
# `journalctl -u sentryusb-telemetry` (NOT `-u sentryusb-tesla-telemetry`)
# to read this daemon's logs.
Description=SentryUSB BLE Telemetry Sampler
Documentation=https://github.com/Sentry-Six/Sentry-USB-Rusty
After=bluetooth.target bluetooth.service dbus.service sentryusb.service
Requires=bluetooth.target
# Wait for the main sentryusb daemon so the drive-data DB has had its
# migrate() run — the sampler relies on the v6 schema being in place
# before it starts inserting into telemetry_samples.
After=sentryusb.service
# Won't start if the BLE binaries haven't been installed yet. The
# sampler binary itself, tesla-control (used to talk to the car), and
# the BLE keypair (set during pairing) all need to be on disk before
# we even attempt to run. Listing each separately means
# `systemctl status` says exactly which one is missing rather than
# spawning and crashing at ExecStart.
ConditionPathExists=/root/bin/sentryusb-tesla-telemetry
ConditionPathExists=/root/bin/tesla-control
ConditionPathExists=/root/.ble/key_private.pem
# Match sentryusb-ble.service's burst limits so a flaky boot doesn't
# burn through the systemd restart budget.
StartLimitBurst=20
StartLimitIntervalSec=120

[Service]
Type=simple
# Clean up our own stale lock at startup (only if we own it AND it's
# >24h old) so a crash mid-session doesn't leave the radio held forever.
ExecStartPre=-/bin/bash -c 'f=/tmp/ble_radio_owner; if [ -e "$f" ] && [ "$(head -n1 "$f" 2>/dev/null)" = "telemetry" ] && [ $(( $(date +%s) - $(stat -c %Y "$f" 2>/dev/null || echo 0) )) -gt 86400 ]; then rm -f "$f"; fi; exit 0'
# Tune bluez LE link parameters for in-car BLE stability.
#
# The Tesla-side BLE link drops every 10-20 min on a moving car with
# the default bluez settings because the supervision_timeout is 720ms.
# A 720ms gap in connection events (RF interference from another
# device, the car body briefly shielding the antenna, a phone-key
# scan burst) triggers a link-dead declaration. After that we have
# to race phone keys for the slot, which we usually lose.
#
# These sysfs writes ask the kernel to use longer-supervision +
# tighter-interval params on future connections, so brief dropouts
# don't tear down the link.
#
# Values are in BLE units:
#   conn_min/max_interval: 1.25ms per unit  (12 = 15ms, 24 = 30ms)
#   supervision_timeout:   10ms per unit    (600 = 6000ms = 6 seconds)
#
# Defaults: min=24/max=40 interval, supervision=72 (720ms)
# New:      min=12/max=24 interval, supervision=600 (6s)
#
# `|| true` everywhere — debugfs isn't always mounted on every Pi
# image, and we'd rather start with default params than block boot.
# Wildcard glob handles hci0/hci1/etc — both onboard and external
# adapters get tuned.
ExecStartPre=-/bin/bash -c 'for d in /sys/kernel/debug/bluetooth/hci*/; do [ -d "$d" ] || continue; echo 600 > "${d}supervision_timeout" 2>/dev/null || true; echo 12 > "${d}conn_min_interval" 2>/dev/null || true; echo 24 > "${d}conn_max_interval" 2>/dev/null || true; echo 0 > "${d}conn_latency" 2>/dev/null || true; done; exit 0'
# Source sentryusb.conf so BLE_ENABLED and TESLA_BLE_VIN are visible
# inside the daemon's env. The Rust code re-reads on each iteration,
# but seeding the initial env here makes the systemd status output
# show the active VIN.
#
# We can't use systemd's `EnvironmentFile=` directive: its parser only
# accepts bare `KEY=value` and silently rejects every line starting
# with `export KEY=value`. The user-facing /root/sentryusb.conf is
# bash-style (archiveloop and other legacy shell scripts source it
# directly), so EnvironmentFile= leaves the daemon with an empty env.
#
# Source via bash instead: `set -a` auto-exports every var sourced,
# `source` understands `export` as a no-op when sourcing, and `exec`
# swaps to the binary so systemd's PID tracking still works.
#
# Each file is gated by `[ -f ] && bash -n` so a hand-edited conf with
# a fat-fingered quote can't crash ExecStart and trigger the same
# Restart=always / StartLimitBurst restart-cycling-to-inactive failure
# mode this commit set out to fix. A broken file is skipped (bash -n
# prints the parse error to the journal so the user sees what went
# wrong); the binary still starts and falls back to its own
# parse_file() reader, which is robust to malformed lines.
ExecStart=/bin/bash -c 'set -a; for f in /root/sentryusb.conf /boot/firmware/sentryusb.conf; do [ -f "$f" ] && bash -n "$f" && source "$f"; done; set +a; exec /root/bin/sentryusb-tesla-telemetry'
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal
SyslogIdentifier=sentryusb-telemetry
# Send SIGTERM and give the process up to 8s to release the radio
# lock and restart sentryusb-ble before SIGKILL. The signal handler
# in main.rs does both before returning.
KillSignal=SIGTERM
TimeoutStopSec=8

[Install]
WantedBy=multi-user.target
